# claude-status-notify.service
[Unit]
Description=Check Claude service status

[Service]
Type=oneshot
ExecStart=%h/bin/claude-status-notify.sh
